# contrib/init/kcoind.conf
description "Kcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env KCOIND_BIN="/usr/bin/kcoind"
env KCOIND_USER="kcoin"
env KCOIND_GROUP="kcoin"
env KCOIND_PIDDIR="/var/run/kcoind"
# upstart can't handle variables constructed with other variables
env KCOIND_PIDFILE="/var/run/kcoind/kcoind.pid"
env KCOIND_CONFIGFILE="/etc/kcoin/kcoin.conf"
env KCOIND_DATADIR="/var/lib/kcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # kcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$KCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run kcoind."
        echo "The setting must appear in $KCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Kcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$KCOIND_PIDDIR"
    chmod 0755 "$KCOIND_PIDDIR"
    chown $KCOIND_USER:$KCOIND_GROUP "$KCOIND_PIDDIR"
    chown $KCOIND_USER:$KCOIND_GROUP "$KCOIND_CONFIGFILE"
    chmod 0660 "$KCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$KCOIND_PIDFILE" \
    --chuid $KCOIND_USER:$KCOIND_GROUP \
    --exec "$KCOIND_BIN" \
    -- \
    -pid="$KCOIND_PIDFILE" \
    -conf="$KCOIND_CONFIGFILE" \
    -datadir="$KCOIND_DATADIR" \
    -disablewallet \
    -daemon

